# templates/default/mailcatcher.upstart.conf.erb
#!upstart
description "Mailcatcher SMTP dev tool"

start on startup
stop on shutdown


<%
command = ["mailcatcher"]
command << "--http-ip #{node['mailcatcher']['http-ip']}"
command << "--http-port #{node['mailcatcher']['http-port']}"
command << "--smtp-ip #{node['mailcatcher']['smtp-ip']}" if node['mailcatcher']['smtp-ip']
command << "--smtp-port #{node['mailcatcher']['smtp-port']}"
command = command.join(" ")
%>

script
  export RBENV_ROOT=/opt/rbenv
  export PATH=$PATH:/opt/rbenv/bin:/opt/rbenv/shims:/usr/local/bin
  eval "$(rbenv init - --no-rehash)"
  rbenv shell 2.0.0-p0
  <%= command %>
end script
